# etc/cron.daily/awstats-update
#!/bin/sh
/usr/local/sbin/awstats-update --quiet $@